# d00m_v0/config/systemd/user/doom-battery-monitor.service
[Unit]
Description=DOOM Battery Monitor Check
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/doom-battery-monitor
LogLevelMax=warning
